# scripts/systemd/thesis-matrix@.service
[Unit]
Description=Thesis matrix shard %i (deepseek-v4-flash, all-go)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/thesis
# .env (LLM routing) is loaded by python-dotenv from the working dir.
# RUNS_ROOT is needed by run_matrix_shard.sh to locate run dirs.
# AGENT_TOOL_MODE is passed through to the runner (default: fixed).
Environment=RUNS_ROOT=runs
Environment=AGENT_TOOL_MODE=fixed
# NUM: total number of shards in the fleet (override to match actual fleet size,
# e.g. `systemctl set-environment NUM=10` or edit EnvironmentFile).
# Default is 10; raise/lower to match the droplet count.
Environment=NUM=10
# One pass over this shard's units (skips done via RUN_COMPLETED.json sentinel);
# on exit systemd restarts it -> a resumable continuous loop.
# Args: SHARD NUM_SHARDS CONC
ExecStart=/bin/bash /root/thesis/scripts/run_matrix_shard.sh %i ${NUM} 2
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target
